# cron/cache-cleanup.cron
0 3 * * * /usr/local/bin/cleanup-cache >> /var/log/cache-cleanup.log 2>&1
